# start_rgb.service
[Unit]
Description=Start rgb server

[Service]
Type=oneshot
ExecStart=/bin/python /etc/rgb_http_server/core.py

[Install]
WantedBy=multi-user.target
